# ops/config/systemd/nginx-iris-relay.socket
[Unit]
Description=Socket for iris-relay nginx proxy

[Socket]
ListenStream=0.0.0.0:16648

[Install]
WantedBy=sockets.target
